# contrib/sm.service
[Unit]
Description=Agentic Sessions Manager
After=default.target

[Service]
ExecStart=%h/.local/bin/sm daemon
Restart=on-failure

[Install]
WantedBy=default.target
